# linux/service/aria2.service
[Unit]
Description=Aria2 Service
After=network.target

[Service]
ExecStart=/usr/local/bin/aria2c --enable-rpc --rpc-listen-all --rpc-allow-origin-all --conf-path=/root/.aria2/aria2.conf

[Install]
WantedBy=default.target
